# systemd/ergctl-resume.service
[Unit]
Description=Re-apply ergctl power mode after resume
# Run AFTER device resume / nvidia-resume completes, so 'ergctl apply' (which
# re-removes any reappeared dGPU audio function) sees the post-resume state.
After=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
After=systemd-suspend.service systemd-hibernate.service systemd-hybrid-sleep.service
After=nvidia-resume.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ergctl apply

[Install]
WantedBy=suspend.target hibernate.target hybrid-sleep.target suspend-then-hibernate.target
