# systemd/ergctl.service
[Unit]
Description=ergctl power mode controller (apply battery/AC/turbo state)
After=multi-user.target
# Don't fight tlp's own apply on boot.
After=tlp.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ergctl apply

[Install]
WantedBy=multi-user.target
